# restic/root/etc/crontabs/root
30  *  *  *  *    /usr/local/bin/backup.sh
10  0  *  *  1    /usr/local/bin/forget.sh