# semanticmw/root/etc/s6-overlay/s6-rc.d/init-smw-config/up
/etc/s6-overlay/s6-rc.d/init-smw-config/run